# prosody/install/etc/s6/services/10-prosody/run
#!/usr/bin/with-contenv bash

while [ ! -f /tmp/state/10-prosody-init ]
do
  sleep 1
done

if [ ! -f /tmp/state/10-prosody ]; then
    chown -R prosody /etc/prosody
    chown -R prosody /usr/lib/prosody
    chown -R prosody /var/lib/prosody
    echo '** [prosody] Registering Focus User'
    sudo -u prosody prosodyctl register focus auth.$HOST $JICOFO_USER_PASS
    
    mkdir -p /tmp/state
    echo 'Initialization Complete' >/tmp/state/10-prosody
fi

echo ''
echo '** [prosody] Starting prosody'
exec s6-setuidgid prosody prosodyctl start
